# utilities/mqtt-to-splunk.service
[Unit]
Description=MQTT to Splunk Forwarder for Utilities Monitoring
Documentation=https://github.com/jinjirosan/HomeMatrixBoard
# Service dependencies - must start after these services
After=network-online.target haproxy.service
Wants=network-online.target
# Require HAProxy to be running - if HAProxy stops, this service stops too
Requires=haproxy.service

[Service]
Type=simple
User=YOUR_USER
Group=YOUR_GROUP
WorkingDirectory=/path/to/HomeMatrixBoard/utilities

# IMPORTANT: Use system Python 3, NOT Splunk's internal Python 2.7
# Explicitly specify /usr/bin/python3 to avoid Splunk's Python
ExecStart=/usr/bin/python3 /path/to/HomeMatrixBoard/utilities/mqtt_to_splunk.py

# Ensure system Python environment is used
Environment="PATH=/usr/local/bin:/usr/bin:/bin"
# Unbuffered output for immediate logs
Environment="PYTHONUNBUFFERED=1"

# Restart policy
Restart=always
RestartSec=10

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=mqtt-to-splunk

# Security hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/path/to/HomeMatrixBoard/utilities

[Install]
WantedBy=multi-user.target

